# infrastructure/projctr.service
[Unit]
Description=Projctr - Portfolio Project Idea Generator
After=network.target
# Start only after this path is reachable (NAS). Comment out if you do not use /mnt/nas.
RequiresMountsFor=/mnt/nas

[Service]
Type=simple
# Must exist on the system and own /opt/projctr (deploy uses campbell).
User=campbell
WorkingDirectory=/opt/projctr
ExecStart=/opt/projctr/projctr
Environment=CONFIG_PATH=/opt/projctr/config.toml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
